# alpine/lirc/lircd.initd
#!/sbin/openrc-run

name=lircd
command="/usr/sbin/$name"
command_args=""
pidfile="/var/run/lirc/lircd.pid"
